# lib/crontab
##### VARIABLES #####
	PATH=/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/home/skilbjo/.local/bin:/home/skilbjo/bin
	ACTIVATE="/home/skilbjo/.virtualenvs/etl/bin/activate"
	MAILTO=jskilbeck@yapstone.com

	APP="/home/skilbjo/python/app/etl"

	ANALYTICS="/home/skilbjo/python/app/etl/sql/Analytics/"
	MMF="/home/skilbjo/python/app/etl/sql/MMF/"
	MPR="/home/skilbjo/python/app/etl/sql/MPR/"	
	TOPDATA="/home/skilbjo/python/app/etl/sql/TopData/"
	EXTNRLGTWY="/home/skilbjo/python/app/etl/sql/ExtrnlGtwy/"	


##### Monthly #####
	## Analytics
		25 11 6 * * skilbjo source $ACTIVATE ; cd $APP ; python etl.py -s crostoli -e "${ANALYTICS}analytics.sql" -l "${ANALYTICS}load_analytics.sql" >/dev/null

	## TopData
		25 11 6 * * skilbjo source $ACTIVATE ; cd $APP ; python etl.py -s crostoli -e "${TOPDATA}top_data.sql" -l "${TOPDATA}load_top_data.sql" >/dev/null

	## MPR
		08 17 5 * * skilbjo source $ACTIVATE ; cd $APP ; python etl.py -s crostoli -e "${MPR}top_data.sql" -l "${MPR}load_top_data.sql" >/dev/null

	## MMF
		25 11 5 * * skilbjo source $ACTIVATE ; cd $APP ; python etl.py -s crostoli -e "${MMF}mmf.sql" -l "${MMF}load_mmf.sql" >/dev/null

	## External Gateway
		12 11 6 * * skilbjo source $ACTIVATE ; cd $APP ; python etl.py -s crostoli -e "${EXTNRLGTWY}extrnlgtwy.sql" -l "${EXTNRLGTWY}load_extrnlgtwy.sql" >/dev/null		

