# lib/crontab
##### VARIABLES #####
	PATH=/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/home/skilbjo/.local/bin:/home/skilbjo/bin
	MAILTO=jskilbeck@yapstone.com
	APP="/home/skilbjo/python/app/etl"
	ANALYTICS=$APP"/sql/Analytics/"
	MMF=$APP"/sql/Analytics/"
	TOPDATA=$APP"/sql/TopData/"


##### Monthly #####
	## Analytics
	
	## TopData

	## MPR

	## MMF
		25 15 5 * * skilbjo cd $APP ; python etl.py -e crostoli -e $TOPDATAmmf.py -l $TOPDATAload_mmf.py >/dev/null

